# daemon/someip_dispatcher.socket.in
[Unit]
Description=SomeIP dispatcher daemon

[Socket]
ListenStream=/tmp/someIPSocket

[Install]
WantedBy=sockets.target
